# deploy/pigenus.service
[Unit]
Description=PiGenus private orchestration node
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pigenus
Group=pigenus
WorkingDirectory=/opt/pigenus
EnvironmentFile=/etc/pigenus/pigenus.env
ExecStart=/opt/pigenus/.venv/bin/uvicorn pigenus.main:app --host ${PIGENUS_HOST} --port ${PIGENUS_PORT}
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/var/lib/pigenus /var/log/pigenus

[Install]
WantedBy=multi-user.target

